# Code/montana-vpn/config-template/xray.service.d/10-donot_touch_single_conf.conf
[Service]
ExecStart=
ExecStart=/usr/local/bin/xray run -config /usr/local/etc/xray/config.json
